# deploy/systemd/finwealth-server.service
[Unit]
Description=Finwealth local ledger API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=finwealth
Group=finwealth
EnvironmentFile=/etc/finwealth/server.env
ExecStart=/opt/finwealth/finwealth-server --ledger-path /var/lib/finwealth/ledger.json
Restart=on-failure
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/finwealth
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
